# profile.d/rvm.sh
#!/bin/sh

[[ -s "/home/awn/.rvm/scripts/rvm" ]] && source "/home/awn/.rvm/scripts/rvm" 
